# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/code/uwsgi.ini
stopsignal=QUIT
stopasagroup=true
startretries=1000000

[program:nginx-app]
command = /usr/sbin/nginx
